# config/postfix/supervisor/conf.d/postfix.ini
[program:postfix]
process_name                    = master
command                         = /usr/lib/postfix/sbin/master -d -c /etc/postfix -s
startsecs                       = 0
autorestart                     = true
stdout_logfile                  = /dev/fd/1
stdout_logfile_maxbytes         = 0
exitcodes                       = 0
