# config/postfix/supervisor/conf.d/rsyslog.ini
[program:rsyslog]
command                         = /usr/sbin/rsyslogd -n
startsecs                       = 5
autorestart                     = true
stdout_logfile                  = /dev/fd/1
stdout_logfile_maxbytes         = 0